# scripts/check-website.service
[Unit]
Description=Check website status

[Service]
ExecStart=/path/to/your/script.sh
